# supervisor/action-tunnel-server.ini
[program:action-tunnel-server]
process_name=%(program_name)s
numprocs=1
command=node index.js
directory=/home/ec2-user/action-tunnel-server
user=ec2-user
environment=HOME="/home/ec2-user", USER="ec2-user"
autostart=true
autorestart=true
stdout_logfile=/home/ec2-user/action-tunnel-server_out.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stderr_logfile=/home/ec2-user/action-tunnel-server_err.log
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=10
